# victoriametrics.service
[Unit]
Description=VictoriaMetrics - time series database
After=network.target

[Service]
User=victoriametrics
Group=victoriametrics
Type=simple
ExecStart=/usr/local/bin/victoria-metrics   -storageDataPath=/var/lib/victoria-metrics   -retentionPeriod=2w   -httpListenAddr=:8428
Restart=always

[Install]
WantedBy=multi-user.target
